# recipes-core/volatile-binds/files/volatile-binds.service.in
[Unit]
Description=Bind mount volatile @where@
DefaultDependencies=false
Before=local-fs.target
After=find_partitions.service
RequiresMountsFor=@whereparent@
ConditionPathIsReadWrite=@whatparent@
ConditionPathExists=@where@
ConditionPathIsReadWrite=!@where@

[Service]
Type=oneshot
RemainAfterExit=Yes
StandardOutput=syslog
TimeoutSec=0
ExecStart=/sbin/mount-copybind @what@ @where@
ExecStop=/sbin/umount-copybind @where@

[Install]
Requires=find_partitions.service
WantedBy=local-fs.target
